# deploy/pipeline-alert.service
[Unit]
Description=Pipeline Alert - Bitbucket pipeline watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Deploy location (repo checked out here, owned by ubuntu)
WorkingDirectory=/home/ubuntu/pipeline-alert

# dotenv ("dotenv/config" in src/config/index.ts) loads .env from WorkingDirectory,
# so no EnvironmentFile needed. Just make sure /home/ubuntu/pipeline-alert/.env exists.

# --- NVM ---
# systemd does NOT source your shell profile, so `nvm` and the `node` on your
# interactive PATH are invisible here. You must hardcode the ABSOLUTE path to the
# nvm-managed node binary. Get it on the box with:
#     nvm which default        (or: nvm which 22 / readlink -f "$(nvm which default)")
# It looks like: /home/ubuntu/.nvm/versions/node/v22.11.0/bin/node
# Replace BOTH the version in Environment=PATH and in ExecStart below to match.
Environment=PATH=/home/ubuntu/.nvm/versions/node/v22.11.0/bin:/usr/local/bin:/usr/bin:/bin

# IMPORTANT: also verify the compiled entry path on the box first:
#   find dist -name index.js
# With the current tsconfig (rootDir "./") it builds to dist/src/index.js.
ExecStart=/home/ubuntu/.nvm/versions/node/v22.11.0/bin/node dist/src/index.js

# Graceful stop: index.ts traps SIGTERM, stops the Telegram long-poll, then exits.
# Giving it time avoids the 409 getUpdates conflict on restart.
KillSignal=SIGTERM
TimeoutStopSec=15

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pipeline-alert

[Install]
WantedBy=multi-user.target
